# roles/backup_mysql/files/etc/cron.d/mysql-backup
0 3 * * * root /usr/local/scripts/mysql-backup.sh
