# roles/piwik_watchdog/files/etc/cron.d/piwik-watchdog
MAILTO=sysops@nonutc.fr
* * * * * root /etc/cron.d/scripts/piwik_watchdog.sh
