# roles/backup_ldap/files/etc/cron.d/ldap-backup
0 4 * * * root /usr/local/scripts/ldap-backup.sh
