# roles/uk2_slavesync/files/etc/cron.d/rsync-slaves
MAILTO=sysops@nonutc.fr
0 * * * * www-data rsync -r rsync-slaves@server50881:/home/www/* /home/www/
