# roles/backup_trac/files/etc/cron.d/backups-trac
4 4 * * * root /usr/local/scripts/dump_trac.sh > /dev/null
