# roles/mail_nonutc_users/files/etc/cron.d/spam-bayes-learn
0 2 * * * root /etc/cron.d/scripts/sa-learn
